# observer-rss.service
[Unit]
Description=Generate Observer RSS

[Service]
Type=oneshot
ExecStart=/home/pbui/src/personal/observer-rss/observer-rss.run
